# jobs/supercmd/com.macshara.supercmd.clear-clipboard-history.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.macshara.supercmd.clear-clipboard-history</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/iamppborah/Developer/personal/macshara/supercmd/clear_clipboard_history.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>11</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/iamppborah/Developer/personal/macshara</string>

  <key>StandardOutPath</key>
  <string>/Users/iamppborah/Developer/personal/macshara/logs/supercmd/supercmd-clipboard-history.launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/iamppborah/Developer/personal/macshara/logs/supercmd/supercmd-clipboard-history.launchd.err.log</string>
</dict>
</plist>
